# infra/docker/certbot/crontab
30 2 * * * /usr/bin/certbot renew --renew-hook /concat-certs.sh
