# infra/docker/certbot/crontab
30 2 * * * /usr/local/bin/certbot renew --renew-hook /renew-hook.sh
